# data/azone.cron
# /etc/cron.d/azone
#
# Weekly DNSSEC zone re-signing job generated by azone

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/etc/nsd:/etc/nsd/azone

1 3 * * 0 root azone sign example.com && nsd-control reload
